# new/deploy/portfolio-new.service
[Unit]
Description=Aravinthan portfolio (new Next.js app)
After=network.target
Wants=network-online.target

[Service]
Type=simple
# setup-vm.sh rewrites these defaults when EC2_USER differs.
User=ec2-user
Group=ec2-user
WorkingDirectory=/home/ec2-user/new
Environment=NODE_ENV=production
EnvironmentFile=/home/ec2-user/.env
ExecStartPre=/bin/mkdir -p /home/ec2-user/new/public/uploads/blogs
ExecStart=/usr/bin/npm run start
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
